# rootfs-scripts/razer-panel-idle-blank.service
[Unit]
Description=Blank the panel and turn off its backlight
After=multi-user.target
ConditionKernelCommandLine=razer_panel_idle_blank

[Service]
Type=oneshot
Environment=RAZER_PANEL_IDLE_SECONDS=60
ExecStart=/usr/local/sbin/razer-panel-idle-blank

[Install]
WantedBy=multi-user.target
