# rootfs-scripts/razer-charge-limits.service
[Unit]
Description=Prefer external power and keep Razer Phone 2 battery between 40 and 80 percent
After=multi-user.target

[Service]
Type=simple
Environment=RAZER_CHARGE_START_PERCENT=40
Environment=RAZER_CHARGE_STOP_PERCENT=80
Environment=RAZER_CHARGE_POLL_SECONDS=30
ExecStart=/usr/local/sbin/razer-charge-limits
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
